# Ubuntu/track.service
# /etc/systemd/system/your_program.service
[Unit]
Description=Your Program with Auto-Restart
After=network.target

[Service]
Type=simple
ExecStart=/path/to/your/program
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5
User=yourusername
WorkingDirectory=/path/to/your/program/directory

# 当连续5次在60秒内启动失败时，执行系统重启
StartLimitAction=reboot-force

[Install]
WantedBy=multi-user.target